# systemd/.config/systemd/user/mega-cmd-server.service
[Unit]
Description=MegaCMD Server
After=network.target

[Service]
ExecStart=/usr/bin/mega-cmd-server

[Install]
WantedBy=default.target
