# systemd/.config/systemd/user/numlockx.service
[Unit]
Description=Numlock Daemon

[Service]
ExecStart=/usr/bin/numlockx

[Install]
WantedBy=default.target
